# ops/systemd/tradeo-postsession-improvement-agent.service
[Unit]
Description=Tradeo post-session improvement agent
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/vboxuser/tradeo
Environment=PYTHONPATH=/home/vboxuser/tradeo/backend
ExecStart=/usr/bin/env python3 /home/vboxuser/tradeo/scripts/run_postsession_improvement_agent.py

