# ops/systemd/user/tradeo-lab-automation-001@.service
[Unit]
Description=Tradeo Lab Automation 001 %i

[Service]
Type=oneshot
WorkingDirectory=/tmp/tradeo-lab-foxhunter-001
Environment=TRADEO_LAB_AUTOMATION_ONLY=1
Environment=TRADEO_NO_LIVE=1
ExecStart=/tmp/tradeo-lab-foxhunter-001/scripts/run_lab_automation_001_daily.sh %i
Nice=5
