# ops/systemd/user/tradeo-lab-daily-mid_collector.service
[Unit]
Description=Tradeo Lab daily mid_collector

[Service]
Type=oneshot
WorkingDirectory=/tmp/tradeo-lab-foxhunter-001
Environment=TRADEO_LAB_AUTOMATION_ONLY=1
Environment=TRADEO_NO_LIVE=1
ExecStart=/tmp/tradeo-lab-foxhunter-001/.venv/bin/python /tmp/tradeo-lab-foxhunter-001/scripts/build_lab_nightly_report.py --phase mid-session
