# ops/systemd/user/tradeo-lab-daily-session.service
[Unit]
Description=Tradeo Lab daily session

[Service]
Type=oneshot
WorkingDirectory=/tmp/tradeo-lab-foxhunter-001
Environment=TRADEO_LAB_AUTOMATION_ONLY=1
Environment=TRADEO_NO_LIVE=1
ExecStart=/bin/bash /tmp/tradeo-lab-foxhunter-001/scripts/run_lab_daily_session.sh
